# debian/cron.d
#  Restart saslauthd if it has died
* * * * * root /usr/lib/vpostmaster/bin/vpm-restartsasl
