# debian/cron.hourly
#!/bin/bash
[ -x /usr/lib/vpostmaster/bin/vpm-pgmaintain ] && \
		/bin/su postgres -c /usr/lib/vpostmaster/bin/vpm-pgmaintain
[ -x /usr/lib/vpostmaster/bin/vpm-backup ] && \
		/usr/lib/vpostmaster/bin/vpm-backup
